# files/etc/crontabs/root
30 5 * * * sleep 90 && touch /etc/banner && wifi
