# files/etc/uci-defaults/101_delete-wan-firewall
#!/bin/sh

if [ "$(uci -q get xsdk.@factory[0].init)" != "1" ]
then
	if [ "$(uci -q get firewall.@zone[1].name)" = "wan" ]
	then
		uci del firewall.@zone[1]
		uci commit firewall
	fi
fi

exit 0
