# files/etc/uci-defaults/102_set-wan-dhcp
#!/bin/sh

if [ "$(uci -q get xsdk.@factory[0].init)" != "1" ]
then
    # set wan dhcp.
    uci set network.wan.proto='dhcp'
    uci commit network
fi

exit 0
