# files/etc/uci-defaults/104_create-all-pre-networks
#!/bin/sh

router_model=$(ubus call system board | jsonfilter -e "@['board_name']")


if [ "$(uci -q get xsdk.@factory[0].init)" != "1" -a "$router_model" = "glinet,gl-mt6000" -o "$router_model" = "zyxel,ex5601-t0-ubootmod" ]
then


uci batch << EOI

add_list network.@device[0].ports=$([[ "$router_model" = "glinet,gl-mt6000" ]] && echo "eth1" || echo "eth0")
set network.@device[0].bridge_empty='1'

add network bridge-vlan
set network.@bridge-vlan[-1].device='br-lan'
set network.@bridge-vlan[-1].vlan='1'
add_list network.@bridge-vlan[-1].ports=$([[ "$router_model" = "glinet,gl-mt6000" ]] && echo "eth1:t" || echo "eth0:t")

add network bridge-vlan
set network.@bridge-vlan[-1].device='br-lan'
set network.@bridge-vlan[-1].vlan='50'
add_list network.@bridge-vlan[-1].ports=$([[ "$router_model" = "glinet,gl-mt6000" ]] && echo "eth1:t" || echo "eth0:t")

add network bridge-vlan
set network.@bridge-vlan[-1].device='br-lan'
set network.@bridge-vlan[-1].vlan='51'
add_list network.@bridge-vlan[-1].ports=$([[ "$router_model" = "glinet,gl-mt6000" ]] && echo "eth1:t" || echo "eth0:t")

add network bridge-vlan
set network.@bridge-vlan[-1].device='br-lan'
set network.@bridge-vlan[-1].vlan='52'
add_list network.@bridge-vlan[-1].ports=$([[ "$router_model" = "glinet,gl-mt6000" ]] && echo "eth1:t" || echo "eth0:t")

add network bridge-vlan
set network.@bridge-vlan[-1].device='br-lan'
set network.@bridge-vlan[-1].vlan='90'
add_list network.@bridge-vlan[-1].ports=$([[ "$router_model" = "glinet,gl-mt6000" ]] && echo "eth1:t" || echo "eth0:t")

add network bridge-vlan
set network.@bridge-vlan[-1].device='br-lan'
set network.@bridge-vlan[-1].vlan='62'
add_list network.@bridge-vlan[-1].ports=$([[ "$router_model" = "glinet,gl-mt6000" ]] && echo "eth1:t" || echo "eth0:t")

add network bridge-vlan
set network.@bridge-vlan[-1].device='br-lan'
set network.@bridge-vlan[-1].vlan='178'
add_list network.@bridge-vlan[-1].ports=$([[ "$router_model" = "glinet,gl-mt6000" ]] && echo "eth1:t" || echo "eth0:t")

add network bridge-vlan
set network.@bridge-vlan[-1].device='br-lan'
set network.@bridge-vlan[-1].vlan='179'
add_list network.@bridge-vlan[-1].ports=$([[ "$router_model" = "glinet,gl-mt6000" ]] && echo "eth1:t" || echo "eth0:t")

add network bridge-vlan
set network.@bridge-vlan[-1].device='br-lan'
set network.@bridge-vlan[-1].vlan='224'
add_list network.@bridge-vlan[-1].ports="lan1:u*"
add_list network.@bridge-vlan[-1].ports="lan2:u*"
add_list network.@bridge-vlan[-1].ports="lan3:u*"
add_list network.@bridge-vlan[-1].ports="lan4:u*"
if [ "$router_model" = "glinet,gl-mt6000" ]; then
	add_list network.@bridge-vlan[-1].ports="lan5:u*"
fi


set network.maintainence='interface'
set network.maintainence.proto='static'
set network.maintainence.ipaddr='192.168.85.1'
set network.maintainence.netmask='255.255.255.0'
set network.maintainence.device='br-lan.244'
set network.maintainence.classlessroute='1'
set network.maintainence.defaultroute='0'

set dhcp.maintainence='dhcp'
set dhcp.maintainence.interface='maintainence'
set dhcp.maintainence.start='2'
set dhcp.maintainence.limit='150'
set dhcp.maintainence.leasetime='12h'

set dhcp.ayaneo='dhcp'
set dhcp.ayaneo.interface='ayaneo'
set dhcp.ayaneo.ignore='1'


set network.ayaneo='interface'
set network.ayaneo.proto='dhcp'
set network.ayaneo.device='br-lan.90'
set network.ayaneo.classlessroute='1'
set network.ayaneo.defaultroute='0'

set dhcp.ayaneo='dhcp'
set dhcp.ayaneo.interface='ayaneo'
set dhcp.ayaneo.ignore='1'

set network.iot='interface'
set network.iot.proto='dhcp'
set network.iot.device='br-lan.52'
set network.iot.classlessroute='1'
set network.iot.defaultroute='0'

set dhcp.iot='dhcp'
set dhcp.iot.interface='iot'
set dhcp.iot.ignore='1'

set network.aqara='interface'
set network.aqara.proto='dhcp'
set network.aqara.device='br-lan.178'
set network.aqara.classlessroute='1'
set network.aqara.defaultroute='0'

set dhcp.aqara='dhcp'
set dhcp.aqara.interface='aqara'
set dhcp.aqara.ignore='1'

set network.hwnet='interface'
set network.hwnet.proto='dhcp'
set network.hwnet.device='br-lan.179'
set network.hwnet.classlessroute='1'
set network.hwnet.defaultroute='0'

set dhcp.hwnet='dhcp'
set dhcp.hwnet.interface='hwnet'
set dhcp.hwnet.ignore='1'

set network.wlan0='interface'
set network.wlan0.proto='dhcp'
set network.wlan0.device='br-lan.50'
set network.wlan0.classlessroute='1'
set network.wlan0.defaultroute='0'

set dhcp.wlan0='dhcp'
set dhcp.wlan0.interface='wlan0'
set dhcp.wlan0.ignore='1'

set network.wlan1='interface'
set network.wlan1.proto='dhcp'
set network.wlan1.device='br-lan.51'
set network.wlan1.classlessroute='1'
set network.wlan1.defaultroute='0'

set dhcp.wlan1='dhcp'
set dhcp.wlan1.interface='wlan1'
set dhcp.wlan1.ignore='1'

set network.betatest='interface'
set network.betatest.proto='dhcp'
set network.betatest.device='br-lan.62'
set network.betatest.classlessroute='1'
set network.betatest.defaultroute='0'

set dhcp.betatest='dhcp'
set dhcp.betatest.interface='betatest'
set dhcp.betatest.ignore='1'

EOI

uci commit network
uci commit dhcp

fi

exit 0
